# deploy/systemd/canvas-agent-digest-weekly.timer
[Unit]
Description=Canvas Course Agent - Weekly digest schedule (Sun 00:00 America/Vancouver)

[Timer]
# systemd timer OnCalendar uses system timezone. We set the *server* timezone to America/Vancouver.
OnCalendar=Sun *-*-* 00:00:00
Persistent=true
Unit=canvas-agent-digest.service

[Install]
WantedBy=timers.target
